# recipes-connectivity/a2dp-agent/files/a2dp-agent.sh
#!/bin/sh
### BEGIN INIT INFO
# Provides:          a2dp-agent
# Required-Start:
# Required-Stop:
# Default-Start:     5
# Default-Stop:
# Short-Description: Pairing agent for A2DP Bluetooth service
### END INIT INFO

/usr/bin/a2dp-agent &

